# etcd/configuration/etcd.service
[Unit]
Description=etcd distributed reliable key-value store

[Service]
Type=notify
ExecStart=/usr/local/bin/etcd
Restart=always